# configs/systemd/llm-distributed.service
# Example — generated by: llm-server.sh dist-install
# Do not copy manually; the script creates this with correct paths, IPs, and model settings.
[Unit]
Description=Distributed LLM Server [vulkan] (2-node RPC via USB4 v2)
After=rpc-server.service
Requires=rpc-server.service

[Service]
Type=simple
Environment=LD_LIBRARY_PATH=/home/<user>/llama.cpp/build/bin
Environment=GGML_VK_VISIBLE_DEVICES=
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do nc -z <worker-tb-ip> 50052 && exit 0; sleep 2; done; echo "Remote RPC not reachable"; exit 1'
ExecStart=/home/<user>/llama.cpp/build/bin/llama-server -m /home/<user>/models/gguf/UD-Q6_K_XL/MiniMax-M2.5-UD-Q6_K_XL-00001-of-00005.gguf --rpc <head-tb-ip>:50052,<worker-tb-ip>:50052 -ngl 999 -ts 1,1 -c 131072 -np 1 -b 4096 -fa 1 --no-mmap -ctk q8_0 -ctv q8_0 --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
