# configs/systemd/rpc-server.service
# Example — generated by: llm-server.sh dist-install
# Do not copy manually; the script creates this with correct paths and IPs.
[Unit]
Description=llama.cpp RPC Server [vulkan] (exposes local GPU for distributed inference)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=LD_LIBRARY_PATH=/home/<user>/llama.cpp/build/bin
# Vulkan RPC: no extra env needed
ExecStart=/home/<user>/llama.cpp/build/bin/rpc-server -H <thunderbolt-ip> -p 50052
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
